# MazeServer2.0.0/conf/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;password=123               ; (default is no password (open server))

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9011        ; (ip_address:port specifier, *:port for all iface)
username=xiaoyao              ; (default is no username (open server))
password=xiaoyao               ; (default is no password (open server))

[supervisord]
logfile=/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
environment=GODEBUG="gctrace=1"
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)
; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
serverurl=http://127.0.0.1:9011 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

; business msg center
[program:gnatsd]
command=/server/bin/gnatsd -c /server/bin/gnatsd.conf
autorestart=true
autostart=true
stdout_logfile=/logs/gnats/stdout_gnatsd.log
stderr_logfile=/logs/gnats/stderr_gnatsd.log

; alert msg center
[program:gnatsd0]
command=/server/bin/gnatsd -c /server/bin/gnatsd0.conf
autorestart=true
autostart=true
stdout_logfile=/logs/gnats/stdout_gnatsd0.log
stderr_logfile=/logs/gnats/stderr_gnatsd0.log

; apns msg center
[program:gnatsd1]
command=/server/bin/gnatsd -c /server/bin/gnatsd1.conf
autorestart=true
autostart=true
stdout_logfile=/logs/gnats/stdout_gnatsd1.log
stderr_logfile=/logs/gnats/stderr_gnatsd1.log

; main business gateway service
[program:main.battery_gateway_server]
command=/server/bin/latest/battery_gateway_server -config=/server/bin/battery_gateway_server.ini
autorestart=true
autostart=false
stdout_logfile=/logs/gateway/main.stdout_gateway.log
stderr_logfile=/logs/gateway/main.stderr_gateway.log

; business alert
[program:main.battery_mail_server]
command=/server/bin/latest/battery_mail_server -config=/server/bin/battery_mail_server.ini
autorestart=true
autostart=false
stdout_logfile=/logs/mail/main.stdout_mail.log
stderr_logfile=/logs/mail/main.stderr_mail.log

; user maintenance service
[program:main.battery_maintenance_server]
command=/server/bin/latest/battery_maintenance_server -config=/server/bin/battery_maintenance_server.ini
autorestart=true
autostart=false
stdout_logfile=/logs/maintenance/main.stdout_maintenance.log
stderr_logfile=/logs/maintenance/main.stderr_maintenance.log

; sina weibo authority service 
[program:main.battery_file_server]
command=/server/bin/latest/battery_file_server -config=/server/bin/battery_file_server.ini
autorestart=true
autostart=false
stdout_logfile=/logs/file/main.stdout_file.log
stderr_logfile=/logs/file/main.stderr_file.log

; sina weibo relocate service
[program:main.battery_file_server_0]
command=/server/bin/latest/battery_file_server -config=/server/bin/battery_file_server_0.ini
autorestart=true
autostart=false
stdout_logfile=/logs/file/main.stdout_file0.log
stderr_logfile=/logs/file/main.stderr_file0.log

; apns provider
[program:main.battery_apns_server]
command=/server/bin/latest/battery_apns_server -config=/server/bin/battery_apns_server.ini
autorestart=true
autostart=false
stdout_logfile=/logs/apns/main.stdout_apns.log
stderr_logfile=/logs/apns/main.stderr_apns.log

[program:primary.battery_app_server_00]
command=/server/bin/latest/battery_app_server -config=/server/bin/battery_app_server.ini -nodeid=0
autorestart=true
autostart=false
stdout_logfile=/logs/app/main.stdout_app_00.log
stderr_logfile=/logs/app/main.stderr_app_00.log

[program:primary.battery_app_server_01]
command=/server/bin/latest/battery_app_server -config=/server/bin/battery_app_server.ini -nodeid=1
autorestart=true
autostart=false
stdout_logfile=/logs/app/main.stdout_app_01.log
stderr_logfile=/logs/app/main.stderr_app_01.log

[program:primary.battery_app_server_02]
command=/server/bin/latest/battery_app_server -config=/server/bin/battery_app_server.ini -nodeid=2
autorestart=true
autostart=false
stdout_logfile=/logs/app/main.stdout_app_02.log
stderr_logfile=/logs/app/main.stderr_app_02.log

[program:primary.battery_app_server_03]
command=/server/bin/latest/battery_app_server -config=/server/bin/battery_app_server.ini -nodeid=3
autorestart=true
autostart=false
stdout_logfile=/logs/app/main.stdout_app_03.log
stderr_logfile=/logs/app/main.stderr_app_03.log

[program:primary.battery_transaction_server_00]
command=/server/bin/latest/battery_transaction_server -config=/server/bin/battery_transaction_server.ini -nodeid=0
autorestart=true
autostart=false
stdout_logfile=/logs/transaction/main.stdout_transaction_00.log
stderr_logfile=/logs/transaction/main.stderr_transaction_00.log

[program:primary.battery_transaction_server_01]
command=/server/bin/latest/battery_transaction_server -config=/server/bin/battery_transaction_server.ini -nodeid=1
autorestart=true
autostart=false
stdout_logfile=/logs/transaction/main.stdout_transaction_01.log
stderr_logfile=/logs/transaction/main.stderr_transaction_01.log

[program:primary.battery_transaction_server_02]
command=/server/bin/latest/battery_transaction_server -config=/server/bin/battery_transaction_server.ini -nodeid=2
autorestart=true
autostart=false
stdout_logfile=/logs/transaction/main.stdout_transaction_02.log
stderr_logfile=/logs/transaction/main.stderr_transaction_02.log

[program:primary.battery_transaction_server_03]
command=/server/bin/latest/battery_transaction_server -config=/server/bin/battery_transaction_server.ini -nodeid=3
autorestart=true
autostart=false
stdout_logfile=/logs/transaction/main.stdout_transaction_03.log
stderr_logfile=/logs/transaction/main.stderr_transaction_03.log

[group:primary]
programs=primary.battery_app_server_00,primary.battery_app_server_01,primary.battery_app_server_02,primary.battery_app_server_03,primary.battery_transaction_server_00,primary.battery_transaction_server_01,primary.battery_transaction_server_02,primary.battery_transaction_server_03

